# review.ini
[program:review]
directory=/home/skuriles/review/src
command=node app.js